# systemd/sensor-service.service
[Unit]
Description=Environmental Sensor Data Collection
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/office-space-db
Environment=PYTHONPATH=/home/pi/office-space-db
EnvironmentFile=/home/pi/office-space-db/config/.env
ExecStart=/home/pi/office-space-db/venv/bin/python3 /home/pi/office-space-db/src/sensor_service.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target